# rootfs/etc/systemd/system.control/chromium.service
[Unit]
Requires=novnc.service
After=novnc.service

[Service]
Environment=DISPLAY=:20
ExecStart=/usr/bin/chromium-browser
Restart=on-failure
User=user
Group=user

[Install]
WantedBy=multi-user.target
