# rootfs/etc/systemd/system.control/warsaw-install.service
[Service]
Type=oneshot
ExecStart=/usr/bin/dpkg -i /w.deb
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
